# base/etc/xdg/autostart/aos-live-icons.desktop
[Desktop Entry]
Type=Application
Name=aos-live-icons.desktop
Exec=/etc/alpha-scripts/aos-live-icons.sh
NoDisplay=true
